# systemd/market-report-agent.service
[Unit]
Description=Market Report Agent — continuous scheduler
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/api-report-agent
EnvironmentFile=/opt/api-report-agent/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/api-report-agent/.venv/bin/python scripts/market_report_agent.py
Restart=always
RestartSec=10
User=deploy
Group=deploy
StandardOutput=journal
StandardError=journal
SyslogIdentifier=market-report-agent

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/api-report-agent/data /opt/api-report-agent/reports /opt/api-report-agent/logs /opt/api-report-agent/runtime
ReadOnlyPaths=/opt/api-report-agent

[Install]
WantedBy=multi-user.target
